# systemd/ge_proton_updater.timer
[Unit]
Description=Run ge_proton_updater every 

[Timer]
OnCalendar=*-*-1,15 18:00:00
Persistent=true

[Install]
WantedBy=timers.target